# modules/masterless/templates/puppet-run.service.erb
[Unit]
Description=Run puppet

[Service]
Type=simple
ExecStart=<%= "#{@bindir}/puppet-run" %>

