# modules/masterless/templates/puppet-run.timer.erb
[Unit]
Description=Schedule puppet run

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hour

[Install]
WantedBy=multi-user.target

