# meta/systemd/puppet-apply.service
[Unit]
Description=Run puppet

[Service]
Type=simple
ExecStart=/opt/goblin/meta/run

