# meta/systemd/puppet-apply.service
[Unit]
Description=Run puppet

[Service]
Type=simple
ExecStart=/usr/local/bin/puppet-run

